# playbooks/roles/wireguard/templates/streisand-wireguard.service.j2
[Unit]
Description=streisand-wireguard
After=network.target

[Service]
Type=oneshot
{% for rule in wireguard_firewall_rules %}
ExecStart={{ rule }}
{% endfor %}

[Install]
WantedBy=multi-user.target
